# infra/systemd/firewall.service
[Unit]
Description=
After=network.target
